# noter/daemon/etc/init/noter-daemon.conf
start on (local-filesystems
          and started dbus)

stop on runlevel [06]

expect daemon
respawn

exec noter-daemon
